# thepeoplepitch-updater.service.example
[Unit]
Description=The People's Pitch Updater Service
After=network.target thepeoplepitch.service

[Service]
User=appuser
Group=appuser
WorkingDirectory=/path/to/app
Environment="PATH=/path/to/app/venv/bin"
Environment="APP_URL=http://localhost:5000"
Environment="CRON_API_KEY=your_cron_api_key_here"
ExecStart=/path/to/app/venv/bin/python updater.py
Restart=always

[Install]
WantedBy=multi-user.target
